# docker/container/Metatrader/etc/s6-overlay/s6-rc.d/svc-mt5server/run
#!/usr/bin/with-contenv bash
# s6 service: MT5 Terminal + gRPC Bridge Manager
# =============================================================================
# Manages both MT5 terminal and gRPC bridge server lifecycle.
# Monitors restart token from health_monitor and performs full restart.
#
# Token-based communication:
# - Reads:  /tmp/.mt5-startup-complete (start.sh signals ready)
# - Reads:  /tmp/.mt5-restart-requested (health_monitor signals restart)
# - Clears: /tmp/.mt5-restart-requested (after restart)
# =============================================================================

# =============================================================================
# CONFIGURATION
# =============================================================================
export CONFIG_DIR="${CONFIG_DIR:-/config}"
export WINEPREFIX="${WINEPREFIX:-$CONFIG_DIR/.wine}"
export WINE_PYTHON_PATH="$WINEPREFIX/drive_c/Python/python.exe"
export wine_executable="${wine_executable:-wine}"
export mt5file="$WINEPREFIX/drive_c/Program Files/MetaTrader 5/terminal64.exe"
export mt5server_port="${mt5server_port:-8001}"

# Token files
STARTUP_MARKER="${STARTUP_MARKER:-/tmp/.mt5-startup-complete}"
STARTUP_IN_PROGRESS="${STARTUP_IN_PROGRESS:-/tmp/.mt5-startup-in-progress}"
RESTART_TOKEN="${RESTART_TOKEN:-/tmp/.mt5-restart-requested}"

# Monitoring interval
TOKEN_CHECK_INTERVAL=5

# =============================================================================
# LOGGING
# =============================================================================
log() {
    echo "[svc-mt5server] $*"
}

# =============================================================================
# HELPER FUNCTIONS
# =============================================================================

run_as_abc() {
    s6-setuidgid abc "$@"
}

is_mt5_running() {
    pgrep -f "terminal64.exe" >/dev/null 2>&1
}

is_bridge_running() {
    pgrep -f "wine.*python.*bridge" >/dev/null 2>&1
}

kill_mt5() {
    log "Stopping MT5 terminal..."
    pkill -f "terminal64.exe" 2>/dev/null || true
    sleep 2
    pkill -9 -f "terminal64.exe" 2>/dev/null || true
}

kill_bridge() {
    log "Stopping gRPC bridge server..."
    pkill -f "wine.*python.*bridge" 2>/dev/null || true
    sleep 1
    pkill -9 -f "wine.*python.*bridge" 2>/dev/null || true
}

start_mt5_terminal() {
    if is_mt5_running; then
        log "MT5 terminal already running"
        return 0
    fi

    if [ ! -f "$mt5file" ]; then
        log "WARN: MT5 terminal not installed"
        return 0
    fi

    log "Starting MT5 terminal..."

    MT5_ARGS="/portable"
    MT5_CONFIG_PATH="$WINEPREFIX/drive_c/MT5Config/startup.ini"
    if [ -f "$MT5_CONFIG_PATH" ]; then
        MT5_ARGS="$MT5_ARGS /config:C:\\MT5Config\\startup.ini"
        log "Using config file for auto-login"
    fi

    run_as_abc ${wine_executable} "$mt5file" $MT5_ARGS &

    # Wait for startup (max 30s)
    local waited=0
    while [ $waited -lt 30 ]; do
        if is_mt5_running; then
            log "MT5 terminal started"
            return 0
        fi
        sleep 2
        waited=$((waited + 2))
    done

    log "WARN: MT5 terminal may not have started"
    return 0
}

start_bridge() {
    log "Starting gRPC bridge server on port ${mt5server_port}..."

    bridge_args="--host 0.0.0.0 --port ${mt5server_port}"
    if [ "${MT5_DEBUG:-1}" = "1" ]; then
        bridge_args="${bridge_args} --debug"
    fi

    run_as_abc ${wine_executable} "$WINE_PYTHON_PATH" -m mt5linux.bridge ${bridge_args} &
    BRIDGE_PID=$!

    sleep 3
    if is_bridge_running; then
        log "gRPC bridge server started (PID: $BRIDGE_PID)"
        return 0
    else
        log "ERROR: gRPC bridge server failed to start"
        return 1
    fi
}

check_restart_token() {
    [ -f "$RESTART_TOKEN" ]
}

clear_restart_token() {
    if [ -f "$RESTART_TOKEN" ]; then
        local reason
        reason=$(cat "$RESTART_TOKEN" 2>/dev/null || echo "unknown")
        log "Clearing restart token (reason: $reason)"
        rm -f "$RESTART_TOKEN"
    fi
}

full_restart() {
    log "Performing full restart..."
    kill_bridge
    kill_mt5
    clear_restart_token
    sleep 2
    start_mt5_terminal
    start_bridge
    log "Full restart complete"
}

# =============================================================================
# WAIT FOR DISPLAY
# =============================================================================
log "Waiting for DISPLAY..."
MAX_DISPLAY_WAIT=60
DISPLAY_WAITED=0
while [ -z "${DISPLAY:-}" ] && [ $DISPLAY_WAITED -lt $MAX_DISPLAY_WAIT ]; do
    sleep 2
    DISPLAY_WAITED=$((DISPLAY_WAITED + 2))
    if [ -f /tmp/.X0-lock ] || [ -S /tmp/.X11-unix/X0 ]; then
        export DISPLAY=:0
        break
    fi
done
[ -z "${DISPLAY:-}" ] && export DISPLAY=:0
log "DISPLAY=$DISPLAY"

# =============================================================================
# WAIT FOR STARTUP TOKEN
# =============================================================================
MAX_WAIT=600
WAITED=0
log "Waiting for startup token..."
while [ $WAITED -lt $MAX_WAIT ]; do
    if [ ! -f "$STARTUP_IN_PROGRESS" ] && [ -f "$STARTUP_MARKER" ]; then
        break
    fi
    sleep 5
    WAITED=$((WAITED + 5))
    [ $((WAITED % 30)) -eq 0 ] && log "Still waiting... (${WAITED}s/${MAX_WAIT}s)"
done

if [ -f "$STARTUP_IN_PROGRESS" ] || [ ! -f "$STARTUP_MARKER" ]; then
    log "ERROR: Startup did not complete within ${MAX_WAIT}s"
    exit 1
fi
log "Startup token received"

# =============================================================================
# VERIFY COMPONENTS
# =============================================================================
log "Verifying components..."

if [ ! -f "$WINE_PYTHON_PATH" ]; then
    log "ERROR: Wine Python not found"
    exit 1
fi

if ! run_as_abc ${wine_executable} "$WINE_PYTHON_PATH" -c "import MetaTrader5" 2>/dev/null; then
    log "ERROR: MetaTrader5 pip not installed"
    exit 1
fi
log "Components OK"

# =============================================================================
# INITIAL STARTUP
# =============================================================================
clear_restart_token
start_mt5_terminal
start_bridge

# =============================================================================
# MAIN LOOP: Monitor restart token and process health
# =============================================================================
log "Entering main loop (monitoring restart token)"

while true; do
    sleep "$TOKEN_CHECK_INTERVAL"

    # Check for restart token from health_monitor
    if check_restart_token; then
        log "Restart token detected"
        full_restart
        continue
    fi

    # Check if bridge is still running
    if ! is_bridge_running; then
        log "Bridge process died, restarting..."
        full_restart
        continue
    fi

    # Check if MT5 is still running
    if ! is_mt5_running; then
        log "MT5 process died, restarting..."
        full_restart
        continue
    fi
done
