# docker/container/Metatrader/etc/s6-overlay/s6-rc.d/svc-mt5server/run
#!/usr/bin/with-contenv bash
# s6 service: Wine RPyC server for MT5
#
# This service runs the RPyC server inside Wine, supervised by s6-overlay.
# If the process crashes, s6 will automatically restart it.
# Uses modern MT5Service (NOT deprecated ClassicService/SlaveService).

# Load environment from 00_env.sh (fail-fast if missing or broken)
ENV_FILE="/Metatrader/scripts/00_env.sh"
if [ ! -f "$ENV_FILE" ]; then
    echo "[svc-mt5server] ERROR: Environment file not found: $ENV_FILE"
    exit 1
fi
if ! source "$ENV_FILE"; then
    echo "[svc-mt5server] ERROR: Failed to load environment from $ENV_FILE"
    exit 1
fi

# Verify required variables are set
if [ -z "${wine_executable:-}" ]; then
    echo "[svc-mt5server] ERROR: wine_executable not defined"
    exit 1
fi
if [ -z "${mt5server_port:-}" ]; then
    echo "[svc-mt5server] ERROR: mt5server_port not defined"
    exit 1
fi
if [ -z "${WINE_PYTHON_PATH:-}" ]; then
    echo "[svc-mt5server] ERROR: WINE_PYTHON_PATH not defined"
    exit 1
fi

# Check if Wine prefix was built during Docker build (fast path)
BUILD_MARKER="${WINEPREFIX:-/config/.wine}/.build-complete"

# STARTUP_MARKER and STARTUP_IN_PROGRESS are provided by 00_env.sh (centralized config)

# ============================================================
# Helper: Run command as user abc (owner of Wine prefix)
# ============================================================
run_as_abc() {
    s6-setuidgid abc "$@"
}

# ============================================================
# Helper: Check if Wine Python is working (run as abc)
# ============================================================
check_wine_python() {
    if [ ! -f "$WINE_PYTHON_PATH" ]; then
        echo "[svc-mt5server] Wine Python not found at: $WINE_PYTHON_PATH"
        return 1
    fi
    if ! run_as_abc ${wine_executable} "$WINE_PYTHON_PATH" --version 2>/dev/null; then
        echo "[svc-mt5server] Wine Python not executing properly"
        return 1
    fi
    return 0
}

# ============================================================
# Helper: Check if MetaTrader5 pip package is installed (run as abc)
# ============================================================
check_mt5_pip() {
    if ! run_as_abc ${wine_executable} "$WINE_PYTHON_PATH" -c "import MetaTrader5" 2>/dev/null; then
        echo "[svc-mt5server] MetaTrader5 pip package not installed"
        return 1
    fi
    return 0
}

# ============================================================
# Helper: Check if MT5 terminal is installed
# ============================================================
check_mt5_terminal() {
    if [ ! -f "$mt5file" ]; then
        echo "[svc-mt5server] MT5 terminal not found at: $mt5file"
        return 1
    fi
    return 0
}

# ============================================================
# Helper: Check if MT5 terminal is running
# ============================================================
is_mt5_running() {
    pgrep -f "terminal64.exe" >/dev/null 2>&1
}

# ============================================================
# Helper: Start MT5 terminal if not running (run as abc)
# Uses /config for auto-login (more reliable than /login per MQL5 forum research)
# See: https://www.mql5.com/en/forum/265448
# ============================================================
start_mt5_terminal() {
    if is_mt5_running; then
        echo "[svc-mt5server] MT5 terminal already running"
        return 0
    fi

    if ! check_mt5_terminal; then
        echo "[svc-mt5server] Cannot start MT5 - not installed"
        return 1
    fi

    echo "[svc-mt5server] Starting MT5 terminal..."

    # Check for config file generated by 30_mt5.sh
    # Config is in C:\MT5Config (no spaces) to avoid quoting issues
    MT5_CONFIG_PATH="$WINEPREFIX/drive_c/MT5Config/startup.ini"

    if [ -f "$MT5_CONFIG_PATH" ]; then
        echo "[svc-mt5server] Using config file for auto-login: $MT5_CONFIG_PATH"
        # Path without spaces - no quoting needed
        run_as_abc ${wine_executable} "$mt5file" /portable /config:C:\\MT5Config\\startup.ini &
    else
        echo "[svc-mt5server] No config file found - starting without auto-login"
        run_as_abc ${wine_executable} "$mt5file" /portable &
    fi

    # Wait for terminal to start (max 30s)
    local waited=0
    while [ $waited -lt 30 ]; do
        if is_mt5_running; then
            echo "[svc-mt5server] MT5 terminal started"
            return 0
        fi
        sleep 2
        waited=$((waited + 2))
    done

    echo "[svc-mt5server] WARN: MT5 terminal may not have started"
    return 0  # Don't fail - let RPyC server try anyway
}

# ============================================================
# Wait for DISPLAY to be ready (Wine needs X11)
# ============================================================
echo "[svc-mt5server] Waiting for DISPLAY..."
MAX_DISPLAY_WAIT=60
DISPLAY_WAITED=0
while [ -z "${DISPLAY:-}" ] && [ $DISPLAY_WAITED -lt $MAX_DISPLAY_WAIT ]; do
    sleep 2
    DISPLAY_WAITED=$((DISPLAY_WAITED + 2))
    # Try to get DISPLAY from environment
    if [ -f /tmp/.X0-lock ] || [ -S /tmp/.X11-unix/X0 ]; then
        export DISPLAY=:0
        break
    fi
done

if [ -z "${DISPLAY:-}" ]; then
    # Default to :0 if not set
    export DISPLAY=:0
fi
echo "[svc-mt5server] DISPLAY=$DISPLAY"

# ============================================================
# Wait for startup scripts to complete
# The startup scripts properly initialize Wine prefix and install MT5
# Must wait for:
#   1. No STARTUP_IN_PROGRESS marker (startup is not running)
#   2. STARTUP_MARKER exists (startup completed successfully)
# This prevents race condition with stale markers from previous runs
# ============================================================
MAX_WAIT=600
WAITED=0
echo "[svc-mt5server] Waiting for startup scripts to complete..."
while [ $WAITED -lt $MAX_WAIT ]; do
    # Check both conditions: startup NOT in progress AND startup completed
    if [ ! -f "$STARTUP_IN_PROGRESS" ] && [ -f "$STARTUP_MARKER" ]; then
        break
    fi
    sleep 5
    WAITED=$((WAITED + 5))
    if [ $((WAITED % 30)) -eq 0 ]; then
        if [ -f "$STARTUP_IN_PROGRESS" ]; then
            echo "[svc-mt5server] Still waiting for startup to finish... (${WAITED}s/${MAX_WAIT}s)"
        else
            echo "[svc-mt5server] Waiting for startup marker... (${WAITED}s/${MAX_WAIT}s)"
        fi
    fi
done

if [ -f "$STARTUP_IN_PROGRESS" ] || [ ! -f "$STARTUP_MARKER" ]; then
    echo "[svc-mt5server] ERROR: Startup did not complete within ${MAX_WAIT}s"
    exit 1
fi
echo "[svc-mt5server] Startup complete"

# ============================================================
# Verify components before starting server
# ============================================================
echo "[svc-mt5server] Verifying components..."

# Check Wine Python
if ! check_wine_python; then
    echo "[svc-mt5server] ERROR: Wine Python not working"
    exit 1
fi
echo "[svc-mt5server] Wine Python OK"

# Check MetaTrader5 pip (installed by 30_mt5.sh during startup)
if ! check_mt5_pip; then
    echo "[svc-mt5server] ERROR: MetaTrader5 pip not installed - startup scripts may have failed"
    exit 1
fi
echo "[svc-mt5server] MetaTrader5 pip OK"

# Check and start MT5 terminal
if check_mt5_terminal; then
    start_mt5_terminal
else
    echo "[svc-mt5server] WARN: MT5 terminal not installed - RPyC will start but MT5 functions won't work"
fi

echo "[svc-mt5server] All checks passed, starting RPyC server..."

# Ensure all required directories have correct ownership (run as root before switching to abc user)
echo "[svc-mt5server] Ensuring directory permissions..."
if [ -d "/config" ]; then
    chown -R abc:abc "/config" 2>/dev/null || echo "[svc-mt5server] Warning: Failed to chown /config"
fi
echo "[svc-mt5server] Directory permissions set"

echo "[svc-mt5server] Starting Wine RPyC server on port ${mt5server_port}"
echo "[svc-mt5server] Using mt5linux.bridge module (standalone, no structlog)"

# Build command arguments
bridge_args="--host 0.0.0.0 -p ${mt5server_port}"

# Add debug flag if MT5_DEBUG is enabled (default: 1)
if [ "${MT5_DEBUG:-1}" = "1" ]; then
    bridge_args="${bridge_args} --debug"
    echo "[svc-mt5server] Debug logging enabled"
fi

# Run as user abc (same as KasmVNC) with exec to replace shell
# Uses bridge.py copied from Linux (standalone, uses only logging + rpyc)
exec s6-setuidgid abc \
    ${wine_executable} "$WINE_PYTHON_PATH" -m mt5linux.bridge ${bridge_args}
