# Metatrader/etc/s6-overlay/s6-rc.d/svc-mt5server/run
#!/usr/bin/with-contenv bash
# s6 service: Wine RPyC server for MT5
#
# This service runs the RPyC server inside Wine, supervised by s6-overlay.
# If the process crashes, s6 will automatically restart it.
# Uses modern MT5Service (NOT deprecated ClassicService/SlaveService).

# Load environment from 00_env.sh (fail-fast if missing or broken)
ENV_FILE="/Metatrader/scripts/00_env.sh"
if [ ! -f "$ENV_FILE" ]; then
    echo "[svc-mt5server] ERROR: Environment file not found: $ENV_FILE"
    exit 1
fi
if ! source "$ENV_FILE"; then
    echo "[svc-mt5server] ERROR: Failed to load environment from $ENV_FILE"
    exit 1
fi

# Verify required variables are set
if [ -z "${wine_executable:-}" ]; then
    echo "[svc-mt5server] ERROR: wine_executable not defined"
    exit 1
fi
if [ -z "${mt5server_port:-}" ]; then
    echo "[svc-mt5server] ERROR: mt5server_port not defined"
    exit 1
fi

# Ensure server directory exists
mkdir -p /tmp/mt5linux

# Generate the RPyC server script with modern MT5Service (NO SlaveService/ClassicService)
python3 << 'PYTHON_EOF'
from pathlib import Path
code = '''#!/usr/bin/env python
"""Modern RPyC server with MT5Service (NO STUBS).

Generated by s6 service for execution inside Wine.
Features:
- Modern rpyc.Service (NOT deprecated SlaveService/ClassicService)
- ThreadPoolServer for concurrent handling
- Signal handling (SIGTERM/SIGINT) for clean container stops
- Python 3.13+ type hints
- NO STUBS - fails if MT5 unavailable
"""
from __future__ import annotations

import argparse
import signal
import sys
import threading
from types import FrameType
from typing import Any

import rpyc
from rpyc.utils.server import ThreadPoolServer

# Global server reference for signal handler
_server: ThreadPoolServer | None = None


class MT5Service(rpyc.Service):
    """Modern RPyC service for MT5. NO STUBS."""

    _mt5_module: Any = None
    _mt5_lock = threading.RLock()

    def on_connect(self, conn: rpyc.Connection) -> None:
        with MT5Service._mt5_lock:
            if MT5Service._mt5_module is None:
                import MetaTrader5 as MT5
                MT5Service._mt5_module = MT5
                print("[mt5server] MT5 module loaded")

    def on_disconnect(self, conn: rpyc.Connection) -> None:
        pass

    def exposed_get_mt5(self) -> Any:
        return MT5Service._mt5_module

    def exposed_health_check(self) -> dict[str, Any]:
        return {"healthy": True, "mt5_available": MT5Service._mt5_module is not None}

    def exposed_initialize(
        self,
        path: str | None = None,
        login: int | None = None,
        password: str | None = None,
        server: str | None = None,
        timeout: int | None = None,
        portable: bool = False,
    ) -> bool:
        return MT5Service._mt5_module.initialize(
            path, login, password, server, timeout, portable
        )

    def exposed_login(
        self, login: int, password: str, server: str, timeout: int = 60000
    ) -> bool:
        return MT5Service._mt5_module.login(login, password, server, timeout)

    def exposed_shutdown(self) -> None:
        return MT5Service._mt5_module.shutdown()

    def exposed_version(self) -> tuple[int, int, str] | None:
        return MT5Service._mt5_module.version()

    def exposed_last_error(self) -> tuple[int, str]:
        return MT5Service._mt5_module.last_error()

    def exposed_terminal_info(self) -> Any:
        return MT5Service._mt5_module.terminal_info()

    def exposed_account_info(self) -> Any:
        return MT5Service._mt5_module.account_info()

    def exposed_symbols_total(self) -> int:
        return MT5Service._mt5_module.symbols_total()

    def exposed_symbols_get(self, group: str | None = None) -> Any:
        if group:
            return MT5Service._mt5_module.symbols_get(group=group)
        return MT5Service._mt5_module.symbols_get()

    def exposed_symbol_info(self, symbol: str) -> Any:
        return MT5Service._mt5_module.symbol_info(symbol)

    def exposed_symbol_info_tick(self, symbol: str) -> Any:
        return MT5Service._mt5_module.symbol_info_tick(symbol)

    def exposed_symbol_select(self, symbol: str, enable: bool = True) -> bool:
        return MT5Service._mt5_module.symbol_select(symbol, enable)

    def exposed_copy_rates_from(
        self, symbol: str, timeframe: int, date_from: Any, count: int
    ) -> Any:
        return MT5Service._mt5_module.copy_rates_from(
            symbol, timeframe, date_from, count
        )

    def exposed_copy_rates_from_pos(
        self, symbol: str, timeframe: int, start_pos: int, count: int
    ) -> Any:
        return MT5Service._mt5_module.copy_rates_from_pos(
            symbol, timeframe, start_pos, count
        )

    def exposed_copy_rates_range(
        self, symbol: str, timeframe: int, date_from: Any, date_to: Any
    ) -> Any:
        return MT5Service._mt5_module.copy_rates_range(
            symbol, timeframe, date_from, date_to
        )

    def exposed_copy_ticks_from(
        self, symbol: str, date_from: Any, count: int, flags: int
    ) -> Any:
        return MT5Service._mt5_module.copy_ticks_from(
            symbol, date_from, count, flags
        )

    def exposed_copy_ticks_range(
        self, symbol: str, date_from: Any, date_to: Any, flags: int
    ) -> Any:
        return MT5Service._mt5_module.copy_ticks_range(
            symbol, date_from, date_to, flags
        )

    def exposed_order_calc_margin(
        self, action: int, symbol: str, volume: float, price: float
    ) -> float | None:
        return MT5Service._mt5_module.order_calc_margin(
            action, symbol, volume, price
        )

    def exposed_order_calc_profit(
        self,
        action: int,
        symbol: str,
        volume: float,
        price_open: float,
        price_close: float,
    ) -> float | None:
        return MT5Service._mt5_module.order_calc_profit(
            action, symbol, volume, price_open, price_close
        )

    def exposed_order_check(self, request: dict[str, Any]) -> Any:
        return MT5Service._mt5_module.order_check(request)

    def exposed_order_send(self, request: dict[str, Any]) -> Any:
        return MT5Service._mt5_module.order_send(request)

    def exposed_positions_total(self) -> int:
        return MT5Service._mt5_module.positions_total()

    def exposed_positions_get(
        self,
        symbol: str | None = None,
        group: str | None = None,
        ticket: int | None = None,
    ) -> Any:
        kwargs: dict[str, Any] = {}
        if symbol:
            kwargs["symbol"] = symbol
        if group:
            kwargs["group"] = group
        if ticket:
            kwargs["ticket"] = ticket
        if kwargs:
            return MT5Service._mt5_module.positions_get(**kwargs)
        return MT5Service._mt5_module.positions_get()

    def exposed_orders_total(self) -> int:
        return MT5Service._mt5_module.orders_total()

    def exposed_orders_get(
        self,
        symbol: str | None = None,
        group: str | None = None,
        ticket: int | None = None,
    ) -> Any:
        kwargs: dict[str, Any] = {}
        if symbol:
            kwargs["symbol"] = symbol
        if group:
            kwargs["group"] = group
        if ticket:
            kwargs["ticket"] = ticket
        if kwargs:
            return MT5Service._mt5_module.orders_get(**kwargs)
        return MT5Service._mt5_module.orders_get()

    def exposed_history_orders_total(self, date_from: Any, date_to: Any) -> int | None:
        return MT5Service._mt5_module.history_orders_total(date_from, date_to)

    def exposed_history_orders_get(
        self,
        date_from: Any = None,
        date_to: Any = None,
        group: str | None = None,
        ticket: int | None = None,
        position: int | None = None,
    ) -> Any:
        kwargs: dict[str, Any] = {}
        if date_from:
            kwargs["date_from"] = date_from
        if date_to:
            kwargs["date_to"] = date_to
        if group:
            kwargs["group"] = group
        if ticket:
            kwargs["ticket"] = ticket
        if position:
            kwargs["position"] = position
        if kwargs:
            return MT5Service._mt5_module.history_orders_get(**kwargs)
        return MT5Service._mt5_module.history_orders_get()

    def exposed_history_deals_total(self, date_from: Any, date_to: Any) -> int | None:
        return MT5Service._mt5_module.history_deals_total(date_from, date_to)

    def exposed_history_deals_get(
        self,
        date_from: Any = None,
        date_to: Any = None,
        group: str | None = None,
        ticket: int | None = None,
        position: int | None = None,
    ) -> Any:
        kwargs: dict[str, Any] = {}
        if date_from:
            kwargs["date_from"] = date_from
        if date_to:
            kwargs["date_to"] = date_to
        if group:
            kwargs["group"] = group
        if ticket:
            kwargs["ticket"] = ticket
        if position:
            kwargs["position"] = position
        if kwargs:
            return MT5Service._mt5_module.history_deals_get(**kwargs)
        return MT5Service._mt5_module.history_deals_get()


def graceful_shutdown(signum: int, frame: FrameType | None) -> None:
    """Handle shutdown signals for clean container stops."""
    sig_name = signal.Signals(signum).name
    print(f"[mt5server] Received {sig_name}, shutting down gracefully...")
    if _server is not None:
        _server.close()
    sys.exit(0)


def main() -> int:
    """Run the RPyC server."""
    global _server

    parser = argparse.ArgumentParser(description="Modern RPyC Server for MT5")
    parser.add_argument("--host", default="0.0.0.0", help="Host to bind")
    parser.add_argument("-p", "--port", type=int, default=8001, help="Port")
    args = parser.parse_args()

    # Setup signal handlers for graceful shutdown
    signal.signal(signal.SIGTERM, graceful_shutdown)
    signal.signal(signal.SIGINT, graceful_shutdown)

    print(f"[mt5server] Starting MT5Service on {args.host}:{args.port}")
    print(f"[mt5server] Python {sys.version}")
    print("[mt5server] Using modern rpyc.Service (NO SlaveService/ClassicService)")

    _server = ThreadPoolServer(
        MT5Service,
        hostname=args.host,
        port=args.port,
        reuse_addr=True,
        nbThreads=10,
        protocol_config={"allow_public_attrs": True, "sync_request_timeout": 300},
    )

    try:
        _server.start()
    except KeyboardInterrupt:
        print("[mt5server] Server interrupted")
    except Exception as e:
        print(f"[mt5server] Server error: {e}")
        return 1
    finally:
        if _server is not None:
            _server.close()
        print("[mt5server] Server stopped")

    return 0


if __name__ == "__main__":
    sys.exit(main())
'''
Path('/tmp/mt5linux/server.py').write_text(code)
PYTHON_EOF

echo "[svc-mt5server] Starting Wine RPyC server on port ${mt5server_port}"
echo "[svc-mt5server] Using modern MT5Service (NO deprecated ClassicService/SlaveService)"

# Run as user abc (same as KasmVNC) with exec to replace shell
exec s6-setuidgid abc \
    ${wine_executable} python.exe /tmp/mt5linux/server.py \
    --host 0.0.0.0 -p ${mt5server_port}
