# Metatrader/etc/s6-overlay/s6-rc.d/svc-mt5server/run
#!/usr/bin/with-contenv bash
# s6 service: Wine RPyC server for MT5
#
# This service runs the RPyC server inside Wine, supervised by s6-overlay.
# If the process crashes, s6 will automatically restart it.

# Load environment from 00_env.sh
if [ -f /Metatrader/scripts/00_env.sh ]; then
    source /Metatrader/scripts/00_env.sh
fi

# Default values
mt5server_port=${mt5server_port:-8001}
wine_executable=${wine_executable:-wine}

# Ensure server directory exists
mkdir -p /tmp/mt5linux

# Generate the RPyC server script
python3 << 'PYTHON_EOF'
from pathlib import Path
code = '''#!/usr/bin/env python
"""Classic RPyC server running SlaveService.

Generated by s6 service for execution inside Wine.
"""
import argparse
from rpyc.utils.server import ThreadedServer
from rpyc.core import SlaveService

def main():
    parser = argparse.ArgumentParser(description="RPyC Classic Server")
    parser.add_argument("--host", default="0.0.0.0", help="Host to bind")
    parser.add_argument("-p", "--port", type=int, default=8001, help="Port")
    args = parser.parse_args()

    print(f"[mt5linux] Starting RPyC server on {args.host}:{args.port}")
    server = ThreadedServer(
        SlaveService,
        hostname=args.host,
        port=args.port,
        reuse_addr=True,
    )
    try:
        server.start()
    except KeyboardInterrupt:
        print("[mt5linux] Server stopped")

if __name__ == "__main__":
    main()
'''
Path('/tmp/mt5linux/server.py').write_text(code)
PYTHON_EOF

echo "[svc-mt5server] Starting Wine RPyC server on port ${mt5server_port}"

# Run as user abc (same as KasmVNC) with exec to replace shell
exec s6-setuidgid abc \
    ${wine_executable} python.exe /tmp/mt5linux/server.py \
    --host 0.0.0.0 -p ${mt5server_port}
