# Metatrader/etc/s6-overlay/s6-rc.d/svc-mt5server/run
#!/usr/bin/with-contenv bash
# s6 service: Wine RPyC server for MT5
#
# This service runs the RPyC server inside Wine, supervised by s6-overlay.
# If the process crashes, s6 will automatically restart it.

# Load environment from 00_env.sh (fail-fast if missing or broken)
ENV_FILE="/Metatrader/scripts/00_env.sh"
if [ ! -f "$ENV_FILE" ]; then
    echo "[svc-mt5server] ERROR: Environment file not found: $ENV_FILE"
    exit 1
fi
if ! source "$ENV_FILE"; then
    echo "[svc-mt5server] ERROR: Failed to load environment from $ENV_FILE"
    exit 1
fi

# Verify required variables are set
if [ -z "${wine_executable:-}" ]; then
    echo "[svc-mt5server] ERROR: wine_executable not defined"
    exit 1
fi
if [ -z "${mt5server_port:-}" ]; then
    echo "[svc-mt5server] ERROR: mt5server_port not defined"
    exit 1
fi

# Ensure server directory exists
mkdir -p /tmp/mt5linux

# Generate the RPyC server script with signal handling for graceful shutdown
python3 << 'PYTHON_EOF'
from pathlib import Path
code = '''#!/usr/bin/env python
"""RPyC 6.x server with graceful shutdown support.

Generated by s6 service for execution inside Wine.
Features:
- RPyC 6.x compatible (ThreadedServer + ClassicService)
- Signal handling (SIGTERM/SIGINT) for clean container stops
- Python 3.13+ type hints
- Pydantic 2 configuration validation (optional)
- Structured log output
"""
from __future__ import annotations

import argparse
import signal
import sys
from types import FrameType
from typing import TYPE_CHECKING

from rpyc.core import ClassicService
from rpyc.utils.server import ThreadedServer

if TYPE_CHECKING:
    from collections.abc import Callable

# Optional Pydantic 2 support for configuration validation
try:
    from pydantic import BaseModel, Field

    class ServerConfig(BaseModel):
        """RPyC server configuration with Pydantic 2 validation."""

        host: str = Field(default="0.0.0.0", description="Host to bind")
        port: int = Field(default=8001, ge=1, le=65535, description="Port number")
        reuse_addr: bool = Field(default=True, description="Reuse address")

    HAS_PYDANTIC = True
except ImportError:
    HAS_PYDANTIC = False

    class ServerConfig:  # type: ignore[no-redef]
        """Fallback config without Pydantic."""

        def __init__(self, host: str = "0.0.0.0", port: int = 8001, reuse_addr: bool = True) -> None:
            self.host = host
            self.port = port
            self.reuse_addr = reuse_addr


# Global server reference for signal handler
_server: ThreadedServer | None = None


def graceful_shutdown(signum: int, frame: FrameType | None) -> None:
    """Handle shutdown signals for clean container stops."""
    sig_name = signal.Signals(signum).name
    print(f"[mt5server] Received {sig_name}, shutting down gracefully...")
    if _server is not None:
        _server.close()
    sys.exit(0)


def main() -> int:
    """Run the RPyC server."""
    global _server

    parser = argparse.ArgumentParser(description="RPyC 6.x Server for MT5")
    parser.add_argument("--host", default="0.0.0.0", help="Host to bind")
    parser.add_argument("-p", "--port", type=int, default=8001, help="Port")
    args = parser.parse_args()

    # Create config with validation (Pydantic) or simple dataclass
    config = ServerConfig(host=args.host, port=args.port)

    # Setup signal handlers for graceful shutdown
    signal.signal(signal.SIGTERM, graceful_shutdown)
    signal.signal(signal.SIGINT, graceful_shutdown)

    print(f"[mt5server] Starting RPyC 6.x server on {config.host}:{config.port}")
    print(f"[mt5server] Python {sys.version}")
    print(f"[mt5server] Pydantic config: {HAS_PYDANTIC}")

    _server = ThreadedServer(
        ClassicService,
        hostname=config.host,
        port=config.port,
        reuse_addr=config.reuse_addr,
    )

    try:
        _server.start()
    except KeyboardInterrupt:
        print("[mt5server] Server interrupted")
    finally:
        if _server is not None:
            _server.close()
        print("[mt5server] Server stopped")

    return 0


if __name__ == "__main__":
    sys.exit(main())
'''
Path('/tmp/mt5linux/server.py').write_text(code)
PYTHON_EOF

echo "[svc-mt5server] Starting Wine RPyC server on port ${mt5server_port}"

# Run as user abc (same as KasmVNC) with exec to replace shell
exec s6-setuidgid abc \
    ${wine_executable} python.exe /tmp/mt5linux/server.py \
    --host 0.0.0.0 -p ${mt5server_port}
